# lib/systemd/system/solydxk-system-aptupdate.service
[Unit]
Description=Update Apt Cache
After=network-online.target

[Service]
Type=idle
ExecStart=/usr/bin/apt-get update

[Install]
WantedBy=network-online.target
