# debian/postinst
#!/bin/sh
# postinst script 
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package
#

case "$1" in
    configure|upgrade|update|reconfigure)
    
        if [ -f /etc/default/grub ]; then
	    sed -i -e 's/GRUB_DISTRIBUTOR=.*/GRUB_DISTRIBUTOR=`lsb_release -d -s 2>\/dev\/null || echo SolydXK`/' /etc/default/grub
	fi
	
	if [ -f /etc/lightdm/lightdm.conf ]; then
	    sed -i -e 's/[# ]*greeter-hide-users=.*/greeter-hide-users=false/' /etc/lightdm/lightdm.conf
	fi
        
        CONT=""
        if [ -f /etc/modprobe.d/alsa-base.conf ]; then
            CONT=$(cat /etc/modprobe.d/alsa-base.conf)
        fi
        
        CHK="options snd-hda-intel model"
        case "$CONT" in
            *$CHK* )
                echo "snd-hda-intel model already set"
                ;;
            *) 
                echo "snd-hda-intel model set to auto"
                echo >> /etc/modprobe.d/alsa-base.conf
                echo "options snd-hda-intel model=auto" >> /etc/modprobe.d/alsa-base.conf
        esac
        
        CONT=""
        if [ -f /etc/modprobe.d/alsa-base-blacklist.conf ]; then
            CONT=$(cat /etc/modprobe.d/alsa-base-blacklist.conf)
        fi
        
        CHK="blacklist pcspkr"
        case "$CONT" in
            *$CHK* )
                echo "pcspkr already blacklisted"
                ;;
            *) 
                echo "pcspkr blacklisted"
                echo >> /etc/modprobe.d/alsa-base-blacklist.conf
	        echo "blacklist pcspkr" >> /etc/modprobe.d/alsa-base-blacklist.conf
        esac

	if which glib-compile-schemas >/dev/null 2>&1
	then
	    glib-compile-schemas /usr/share/glib-2.0/schemas
	fi
	
	# Add some settings for Samba
	SMB='/etc/samba/smb.conf'
	if [ -f $SMB ]; then
	  if ! grep -Fq "## SOLYDXK" $SMB; then
	    sed -i -e 's/\[global\]/\[global\]\n  ## SOLYDXK\n\n  follow symlinks = yes\n  wide links = yes\n  unix extensions = no\n  client lanman auth = yes\n  client ntlmv2 auth = no\n/' $SMB
	    sed -i -e 's/^;create mask/create mask/' $SMB
	  fi
	fi
	
	/usr/lib/solydxk/system/adjust.py
	
	# Remove old initv script
	INITCONF='/etc/init.d/solydxk-system'
	if [ -e $INITCONF ]; then
	  update-rc.d solydxk-system remove
	  rm $INITCONF
	fi
	
	systemctl enable solydxk-system
	
	# If dpkg removed /usr/local/bin: recreate it
	if [ ! -e /usr/local/bin ]; then
	  mkdir -p /usr/local/bin
	fi
	
    ;;
    abort-upgrade|abort-remove|abort-deconfigure)

    ;;
    triggered)
	/usr/lib/solydxk/system/adjust.py
    ;;
    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.



exit 0


